# wordpress/etc/services.d/apache2/run
#!/usr/bin/with-contenv sh
set -e

exec 2>&1
exec apachectl -DFOREGROUND
